# services/rplayer.service
[Unit]
Description=RPlayer (Radiko + World Radio)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/pizero2-rplayer
Environment=RPLAYER_ALSA_DEVICE=hw:1,0
Environment=RPLAYER_FONT=/usr/share/fonts/truetype/noto/NotoSansCJK-Regular.ttc
Environment=RPLAYER_ST7789_SPEED_HZ=20000000
ExecStart=/opt/pizero2-rplayer/.venv/bin/python /opt/pizero2-rplayer/rplayer.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
